# cc_usage_reporter/service_templates/linux/cc-usage-reporter.service
[Unit]
Description=CC Usage Reporter Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/cc-usage-reporter
ExecStart=%h/.local/bin/cc-usage-reporter daemon --config %h/.config/cc-switch/usage_reporter.json
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
